# patches/system/to-sort/system-update.patch/files/usr/local/etc/periodic/daily/999.101.system-upgrade
#!/bin/sh

pkg update
pkg upgrade -y

# install any patches
freebsd-update fetch
freebsd-update install
