# patches/system/security/daily-logout-of-sessions.patch/files/usr/local/etc/periodic/daily/999.102.clear-user-sessions
#!/bin/sh

/usr/local/sbin/clear-user-sessions | logger -i -t  "clear-user-sessions (cron @daily)"
